# .github/launch-agents/com.founder.runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.founder.runner</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/carlos.miranda/Developer/founder-os-seed/.github/launch-agents/founder-run.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>14400</integer>
    <key>StandardOutPath</key>
    <string>/Users/carlos.miranda/Developer/founder-os-seed/.github/launch-agents/logs/founder.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/carlos.miranda/Developer/founder-os-seed/.github/launch-agents/logs/founder-error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
