# .github/launch-agents/com.founder.mac-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.founder.mac-runner</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/carlos.miranda/Developer/founder-os-seed/.github/launch-agents/mac-runner-run.sh</string>
    </array>
    <!-- Runs every 4 hours. The run script guards against overlapping with the Founder. -->
    <key>StartInterval</key>
    <integer>14400</integer>
    <key>StandardOutPath</key>
    <string>/Users/carlos.miranda/Developer/founder-os-seed/.github/launch-agents/logs/mac-runner.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/carlos.miranda/Developer/founder-os-seed/.github/launch-agents/logs/mac-runner-error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
